# device/etc/systemd/system/op3t-battery-guard.service
[Unit]
Description=OP3T battery charge-band guard
After=multi-user.target

[Service]
ExecStart=/usr/local/bin/op3t-battery-guard
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
