# build/aports/op3t-helpers/op3t-metrics-web.socket
[Unit]
Description=op3t metrics web viewer socket

[Socket]
ListenStream=8088
Accept=yes

[Install]
WantedBy=sockets.target
