# build/aports/op3t-helpers/op3t-metrics.service
[Unit]
Description=Sample op3t battery/thermal metrics to CSV

[Service]
Type=oneshot
ExecStart=/usr/bin/op3t-metrics-sample
